# supervisord.conf
[supervisord]
nodaemon=true

[program:bench]
command=cd /home/frappe/frappe-bench && bench start 